# scheduler/com.nse.bias-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nse.bias-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shobha/Documents/Python_Learn/.venv/bin/python</string>
        <string>/Users/shobha/Documents/Python_Learn/NSE_Dashboard/scheduler/daily_runner.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/shobha/Documents/Python_Learn/NSE_Dashboard</string>

    <key>StandardOutPath</key>
    <string>/Users/shobha/Documents/Python_Learn/NSE_Dashboard/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shobha/Documents/Python_Learn/NSE_Dashboard/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/shobha/Documents/Python_Learn/NSE_Dashboard</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
